# update-trackers.cron
0 4 * * * /usr/local/bin/update-trackers.sh